# .config/autostart/Solaar.desktop
[Desktop Entry]
Exec=solaar
Icon=
Name=Solaar
Path=
Terminal=True
Type=Application
